# deploy/armn.service
[Unit]
Description=ARM$N Astro/Fastify server
After=network.target

[Service]
Type=simple
User=armsn
WorkingDirectory=/home/armsn/app
EnvironmentFile=/home/armsn/app/.env
Environment="NODE_ENV=production"
Environment="FIRST=false"
ExecStart=/home/armsn/.bun/bin/bun run start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
